# root/etc/services.d/openvpn-monitor/run
#!/usr/bin/with-contenv sh

s6-echo "Starting openvpn-monitor"

cd /openvpn-monitor
gunicorn openvpn-monitor --bind 0.0.0.0:80
